# openclaw/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:bridge]
command=node /app/bridge/server.js
directory=/app/bridge
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

# Uncomment when OpenClaw gateway binary is available:
# [program:openclaw]
# command=/app/openclaw/openclaw --config /app/openclaw.json
# directory=/app/openclaw
# autostart=true
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
